# working/boot-animation.service
[Unit]
Description=CLST2 Boot Animation (Goose)
# Это условие говорит: "Запускайся, только если мы грузимся в текстовую консоль, а не сразу в графику"
# Это полезно для отладки
ConditionPathExists=!/run/systemd/system/display-manager.service

[Service]
Type=simple
# Это главная команда: "Запустить наш скрипт-аниматор и направить его вывод на первый терминал (tty1)"
ExecStart=/bin/sh -c '/usr/local/bin/boot-animation.sh > /dev/tty1'
# Эта опция говорит: "Останови этот сервис, как только он станет не нужен"
# (то есть когда запустится следующий сервис в цепочке - getty)
StopWhenUnneeded=yes

[Install]
# Это самая важная часть. Мы "привязываем" нашу анимацию к сервису getty@tty1.
# getty@tty1 - это сервис, который показывает приглашение "login:".
# Наша анимация будет запускаться ДО него, а останавливаться, когда он запустится.
WantedBy=getty@tty1.service